# deploy/visitantes-cleanup.service.example
[Unit]
Description=Limpeza de visitantes expirados (Sistema de Orcamentos)
After=network.target

[Service]
Type=oneshot
User=www-data
Group=www-data
WorkingDirectory=/var/www/sistema_orcamentos/current
EnvironmentFile=/var/www/sistema_orcamentos/shared/.env
ExecStart=/var/www/sistema_orcamentos/venv/bin/python manage.py limpar_visitantes_expirados
