# deploy/visitantes-cleanup.timer.example
[Unit]
Description=Agenda limpeza de visitantes expirados (a cada 1 hora)

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
